# systemd/sunset.service
[Unit]
Description=Reduces blue light

[Service]
Type=oneshot
ExecStart=hyprsunset -t 3600
