# systemd/sunset.service
[Unit]
Description=Reduces blue light

[Service]
Type=oneshot
ExecStart=/home/mbradley/scripts/sunset.sh
